# x402-server/x402-server.service
[Unit]
Description=x402 Reputation Payment Server
After=network.target

[Service]
Type=simple
User=phan_harry
WorkingDirectory=/home/phan_harry/.openclaw/workspace/x402-server
ExecStart=/usr/bin/node server.mjs
Restart=always
RestartSec=5
EnvironmentFile=/home/phan_harry/.openclaw/.env
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
